# packages/dcos-metrics/extra/dcos-metrics.socket
[Unit]
Description=Metrics Socket: Socket for DC/OS Metrics Service
[Socket]
ListenStream=/run/dcos/dcos-metrics.sock
