# packages/dcos-metrics/extra/dcos-metrics-agent.windows.service
[Unit]
Description=DC/OS Metrics Agent: exposes node, container, and application metrics

[Service]
Restart=always
RestartSec=5
KillSignal=SIGKILL
LimitNOFILE=16384
PermissionsStartOnly=True
StandardOutput=file:c:\var\log\dcos-metrics-agent.log
StandardError=file:c:\var\log\dcos-metrics-agent.log
User=dcos_metrics
SupplementaryGroups=dcos_adminrouter
EnvironmentFile=c:\opt\mesosphere\environment.export.ps1
EnvironmentFile=c:\opt\mesosphere\etc\dcos-metrics.env.ps1
EnvironmentFile=-c:\opt\mesosphere\etc\dcos-metrics-extra.env.ps1
# NOTE: bootstrap here needs access to name resolution, so needs dcos-net working
ExecStartPre=c:\opt\mesosphere\bin\python\scripts\bootstrap.exe dcos-metrics-agent
ExecStartPre=pwsh.exe -command test-path c:\var\lib\dcos\cluster-id
ExecStart=c:\opt\mesosphere\bin\dcos-metrics.exe -role agent -config c:\opt\mesosphere\etc\dcos-metrics-config.yaml
