# packages/dcos-metrics/extra/dcos-metrics-master.service
[Unit]
Description=Metrics Master Service: DC/OS Metrics Master Service for Host Level Metrics.
[Service]
Restart=always
RestartSec=5
PermissionsStartOnly=True
User=dcos_metrics
EnvironmentFile=/opt/mesosphere/environment
ExecStartPre=/opt/mesosphere/bin/bootstrap dcos-metrics-master
ExecStart=/opt/mesosphere/bin/dcos-metrics -config /opt/mesosphere/etc/dcos-metrics-config.yaml -role master
