# packages/dcos-metrics/extra/dcos-metrics-master.socket
[Unit]
Description=Metrics Master Socket: Socket for DC/OS Metrics Master Service
[Socket]
ListenStream=/run/dcos/dcos-metrics-master.sock
