# packages/dcos-metrics/extra/dcos-metrics-agent.socket
[Unit]
Description=Metrics Agent Socket: Socket for DC/OS Metrics Agent Service
[Socket]
ListenStream=/run/dcos/dcos-metrics-agent.sock
