# packages/adminrouter/extra/systemd/dcos-adminrouter-reload.timer
[Unit]
Description=Admin Router Reloader Timer: periodically reloads Admin Router to pick up domain resolution changes
[Timer]
OnBootSec=5sec
OnUnitActiveSec=30s
