# packages/3dt/extra/dcos-3dt-master.service
[Unit]
Description=DC/OS Diagnostics (3DT) Master: aggregates and exposes component health
[Service]
EnvironmentFile=/opt/mesosphere/environment
EnvironmentFile=/opt/mesosphere/etc/dcos-3dt.env
EnvironmentFile=-/opt/mesosphere/etc/dcos-3dt-extra.env
Restart=always
StartLimitInterval=0
RestartSec=5
LimitNOFILE=16384
PermissionsStartOnly=True
ExecStartPre=/opt/mesosphere/bin/bootstrap dcos-3dt-master
User=dcos_3dt
ExecStart=/opt/mesosphere/bin/3dt --3dt-config ${DCOS_3DT_CONFIG_PATH}
