# packages/3dt/extra/dcos-3dt.socket
[Unit]
Description=DC/OS Diagnostics (3DT) Agent Socket: socket for DC/OS Diagnostics Agent
[Socket]
ListenStream=/run/dcos/3dt.sock
