# packages/navstar/extra/dcos-spartan-watchdog.timer
[Unit]
Description=DNS Dispatcher Watchdog Timer: Periodically check is Spartan is working
[Timer]
OnBootSec=5sec
OnUnitActiveSec=5min
