# packages/rexray/extra/dcos-rexray.service
[Unit]
Description=REX-Ray: A vendor agnostic storage orchestration engine

[Service]
StartLimitInterval=0
RestartSec=15
Restart=always
EnvironmentFile=/opt/mesosphere/environment
Environment=REXRAY_HOME=/
ExecStart=__PKG_PATH__/bin/rexray start -f
