# packages/dcos-net/extra/dcos-dns-watchdog.service
[Unit]
Description=DC/OS DNS Watchdog: monitoring the state of DC/OS DNS Forwarder

[Service]
Restart=always
StartLimitInterval=0
RestartSec=15
WorkingDirectory=/opt/mesosphere/active/dcos-net/dcos-net
EnvironmentFile=/opt/mesosphere/environment
Environment=HOME=/opt/mesosphere
ExecStart=/opt/mesosphere/active/dcos-net/bin/dcos-net-observer.py dcos_dns
