# packages/dcos-net/extra/dcos-net-watchdog.timer
[Unit]
Description=DC/OS Net Watchdog Timer: timer to periodically trigger DC/OS Net Watchdog
[Timer]
OnBootSec=5sec
OnUnitActiveSec=2min
