# packages/dcos-history-service/extra/dcos-history.service
[Unit]
Description=Mesos History: DC/OS Resource Metrics History Service/API
After=dcos-mesos-master.service
[Service]
User=dcos_history
Restart=always
StartLimitInterval=0
RestartSec=5
EnvironmentFile=/opt/mesosphere/environment
Environment=HISTORY_BUFFER_DIR=/var/lib/dcos/dcos-history-service
ExecStart=/opt/mesosphere/bin/dcos-history
