# packages/calico/extra/dcos-calico-felix.service
[Unit]
Description=DC/OS calico node: Calico Core component responsible for programming routes and security policy
ConditionPathExists=/opt/mesosphere/etc/calico_enabled

[Service]
EnvironmentFile=/opt/mesosphere/environment
EnvironmentFile=/opt/mesosphere/etc/calico-node.env
ExecStart=/opt/mesosphere/bin/start-calico-felix.sh

Restart=always
StartLimitInterval=0
KillSignal=SIGKILL
RestartSec=5

[Install]
WantedBy=multi-user.target
