# packages/calico/extra/dcos-calico-libnetwork-plugin.service
[Unit]
Description=DC/OS Calico Node: Docker libnetwork plugin for Calico in DC/OS
ConditionPathExists=/opt/mesosphere/etc/calico/calico_enabled

[Service]
Environment=CALICO_LIBNETWORK_LABEL_ENDPOINTS=true
EnvironmentFile=/opt/mesosphere/environment
EnvironmentFile=/opt/mesosphere/etc/calico/calico-node.env
EnvironmentFile=/opt/mesosphere/etc/calico/calico-node-datastore.env
ExecStart=/opt/mesosphere/bin/start-calico-libnetwork-plugin.sh
ExecStartPost=/opt/mesosphere/bin/create-calico-docker-network.py
TimeoutStartSec=180s

Restart=always
StartLimitInterval=0
KillSignal=SIGKILL
RestartSec=5
