# packages/cockroach/extra/dcos-cockroachdb-config-change.timer
[Unit]
Description=DC/OS CockroachDB Config Check Timer
Documentation=https://docs.d2iq.com/mesosphere/dcos/

[Timer]
# initial activation will be done by dcos-cockroach.service
OnUnitActiveSec=10min
