# packages/cockroach/extra/dcos-cockroachdb-config-change.timer
[Unit]
Description=CockroachDB Set Configuration Timer: Timer to periodically set the CockroachDB configuration
Documentation=https://docs.mesosphere.com

[Timer]
OnStartupSec=5min
OnUnitActiveSec=10min
