# packages/cockroach/extra/dcos-cockroachdb-config-change.service
[Unit]
Description=CockroachDB Disable Diagnostics: Disable CockroachDB diagnostic reporting
Documentation=https://docs.mesosphere.com

[Service]
Type=simple
User=dcos_cockroach
StartLimitInterval=0
Restart=on-failure
RestartSec=5
LimitNOFILE=16384
EnvironmentFile=/opt/mesosphere/environment
ExecStart=/opt/mesosphere/active/cockroach/bin/cockroachdb-change-config.py

[Install]
WantedBy=multi-user.target
