# packages/dcos-check-runner/extra/dcos-checks-api.service
[Unit]
Description=DC/OS Checks API: Provides node and cluster checks

[Service]
Restart=always
StartLimitInterval=0

LimitNOFILE=16384

Sockets=dcos-checks-api.socket

PermissionsStartOnly=True
User=dcos_check_runner
# Allow r/w access to the socket file
SupplementaryGroups=dcos_adminrouter
EnvironmentFile=/opt/mesosphere/environment
EnvironmentFile=/opt/mesosphere/etc/dcos-check-runner.env
ExecStart=/opt/mesosphere/bin/dcos-check-runner http-server --config ${DCOS_DIAGNOSTICS_CONFIG_PATH}
