# packages/dcos-check-runner/extra/dcos-checks-api.socket
[Unit]
Description=DC/OS Checks API Socket: socket for DC/OS Checks API

[Socket]
ListenStream=/run/dcos/dcos-checks-api.sock
SocketUser=root
SocketGroup=dcos_adminrouter
SocketMode=0660
