# packages/dcos-check-runner/extra/dcos-checks-poststart.service
[Unit]
Description=DC/OS Poststart Checks: Run node-poststart checks
[Service]
EnvironmentFile=/opt/mesosphere/environment
EnvironmentFile=/opt/mesosphere/etc/dcos-check-runner.env
StandardOutput=journal
User=root
Group=root
ExecStart=/opt/mesosphere/bin/dcos-check-runner check node-poststart --config ${DCOS_CHECK_RUNNER_CONFIG_PATH}
