# packages/spartan/extra/dcos-spartan-watchdog.timer
[Unit]
Description=DNS Dispatcher Watchdog Timer: Periodically check is Spartan is working
[Timer]
OnBootSec=65sec
OnUnitActiveSec=5min
