# packages/spartan/extra/dcos-gen-resolvconf.timer
[Unit]
Description=Generate resolv.conf Timer: periodically updates the network name resolution configuration
[Timer]
OnBootSec=5sec
OnUnitActiveSec=1min
