# packages/spartan/extra/dcos-spartan-watchdog.timer
[Unit]
Description=DNS Forwarder (Spartan) Watchdog Timer: timer to periodically trigger Spartan Watchdog
[Timer]
OnBootSec=5sec
OnUnitActiveSec=2min
