# packages/spartan/extra/dcos-gen-resolvconf.service
[Unit]
Description=Generate resolv.conf: configures network name resolution

[Service]
Type=oneshot
StandardOutput=journal
StandardError=journal
EnvironmentFile=/opt/mesosphere/environment
EnvironmentFile=/opt/mesosphere/etc/dns_config
EnvironmentFile=/opt/mesosphere/etc/dns_search_config
EnvironmentFile=-/opt/mesosphere/etc/dns_config_master
ExecStart=/opt/mesosphere/bin/gen_resolvconf.py /etc/resolv.conf
