# packages/logrotate/extra/dcos-logrotate-agent.service
[Unit]
Description=Logrotate Mesos Agent: Rotate Mesos agent logs
[Service]
Type=simple
EnvironmentFile=/opt/mesosphere/environment
ExecStartPre=/usr/bin/mkdir -p /var/log/mesos/archive
ExecStart=$PKG_PATH/bin/logrotate -v /opt/mesosphere/etc/logrotate_agent.config -s /var/lib/dcos/mesos/logrotate_agent.status
