# packages/logrotate/extra/dcos-logrotate.timer
[Unit]
Description=Logrotate Timer: Timer to trigger every 2 minutes
[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
