# packages/dcos-log/extra/dcos-log-master.service
[Unit]
Description=DC/OS Log Master: exposes master node and component logs

[Service]
EnvironmentFile=/opt/mesosphere/environment
EnvironmentFile=/opt/mesosphere/etc/dcos-log.env
EnvironmentFile=-/opt/mesosphere/etc/dcos-log-extra.env
Restart=always
StartLimitInterval=0
RestartSec=5
LimitNOFILE=16384
User=dcos_log
# Allow r/w access to the socket file
SupplementaryGroups=dcos_adminrouter
ExecStart=/opt/mesosphere/bin/dcos-log -config ${DCOS_LOG_CONFIG_PATH}
Sockets=dcos-log-master.socket
