# packages/dcos-log/extra/dcos-log.socket
[Unit]
Description=DC/OS Log socket: DC/OS logging socket
[Socket]
ListenStream=/run/dcos/dcos-log.sock
