# packages/adminrouter/extra/dcos-adminrouter-reload.timer
[Unit]
Description=Admin Router Reloader Timer: Periodically reload admin router nginx config to pickup new dns
[Timer]
OnBootSec=5sec
OnUnitActiveSec=30s
