# packages/adminrouter/extra/dcos-adminrouter-reload.timer
[Unit]
Description=Admin Router Reloader Timer: periodically reload admin router to pick up domain resolution changes
[Timer]
OnBootSec=5sec
OnUnitActiveSec=30s
