# packages/adminrouter/extra/dcos-adminrouter-reload.service
[Unit]
Description=Admin Router Reloader: Reload admin router to get new DNS

[Service]
Type=oneshot
EnvironmentFile=/opt/mesosphere/environment
ExecStart=-$PKG_PATH/nginx/sbin/adminrouter.sh -c $PKG_PATH/nginx/conf/nginx.master.conf -s reload
