# packages/dcos-ui-service/extra/dcos-ui-service.socket
[Unit]
Description=DC/OS UI Service: service for DC/OS UI
[Socket]
ListenStream=/run/dcos/dcos-ui-service.sock
SocketUser=root
SocketGroup=dcos_adminrouter
SocketMode=0660
