# packages/docker-gc/extra/dcos-docker-gc.service
[Unit]
Description=Docker GC: periodically garbage collects Docker containers and images
ConditionPathExists=/opt/mesosphere/etc/docker_gc_enabled

[Service]
Type=simple
ExecStartPre=-/bin/sh -c 'docker container prune --force'
ExecStartPre=-/bin/sh -c 'docker image prune --force'
ExecStartPre=-/bin/sh -c 'docker volume prune --force'
ExecStart=/bin/true