# packages/docker-gc/extra/dcos-docker-gc.timer
[Unit]
Description=Docker GC Timer: Timer to trigger periodic `docker-gc`
[Timer]
OnBootSec=45sec
OnCalendar=hourly
