# packages/docker-gc/extra/dcos-docker-gc.service
[Unit]
Description=Docker GC: periodically garbage collects Docker containers and images
ConditionPathExists=/opt/mesosphere/etc/docker_gc_enabled

[Service]
Type=simple
ExecStart=/bin/sh -c 'docker system prune --volumes --force'