# packages/docker-gc/extra/dcos-docker-gc.timer
[Unit]
Description=Docker GC Timer: timer to periodically trigger Docker GC
[Timer]
OnBootSec=45sec
OnCalendar=hourly
