# packages/dcos-signal/extra/dcos-signal.timer
[Unit]
Description=DC/OS Signal Timer: timer for DC/OS Signal service
[Timer]
Unit=dcos-signal.service
OnCalendar=hourly
