# packages/dcos-signal/extra/dcos-signal-ee.timer
[Unit]
Description=Signal Timer: Timer for DC/OS Signal Service
[Timer]
Unit=dcos-signal.service
OnCalendar=*:0/5
