# packages/dcos-checks/extra/dcos-checks-poststart.timer
[Unit]
Description=DC/OS Checks Timer: timer for DC/OS Checks service
[Timer]
OnCalendar=minutely
AccuracySec=10
