# packages/telegraf/extra/dcos-telegraf-dcos_metrics.socket
[Unit]
Description=Telegraf dcos_metrics Socket: socket for Telegraf's dcos_metrics plugin
[Socket]
Service=dcos-telegraf.service
ListenStream=/run/dcos/telegraf-dcos-metrics.sock
SocketUser=root
SocketGroup=dcos_adminrouter
SocketMode=0660
