# packages/mesos-dns/extra/dcos-mesos-dns.service
[Unit]
Description=Mesos DNS: DNS based Service Discovery
After=dcos-mesos-master.service
[Service]
User=dcos_mesos_dns
Restart=always
StartLimitInterval=0
RestartSec=5
EnvironmentFile=/opt/mesosphere/environment
ExecStartPre=/opt/mesosphere/bin/bootstrap --ignore-unknown dcos-mesos-dns
ExecStart=/opt/mesosphere/bin/mesos-dns --config=/opt/mesosphere/etc/mesos-dns.json -logtostderr=true
