# packages/bouncer/extra/dcos-bouncer-migrate-users.service
[Unit]
Description=DC/OS Identity and Access dcos-oauth users migration
Documentation=https://docs.d2iq.com/mesosphere/dcos/
# Wait for Bouncer to be available
After=dcos-bouncer.service

[Service]
Type=simple
User=dcos_bouncer
StartLimitInterval=0
Restart=on-failure
RestartSec=30
LimitNOFILE=16384
EnvironmentFile=/opt/mesosphere/environment
ExecStart=/opt/mesosphere/active/bouncer/bin/iam-migrate-users-from-zk.py

[Install]
WantedBy=multi-user.target
