# packages/mesos/extra/dcos-vol-discovery-pub-agent.service
[Unit]
Description=Mesos Persistent Volume Discovery: Discovers /dcos/volume[n] and generates MESOS_RESOURCES
ConditionPathExists=!/var/lib/dcos/mesos-resources

[Service]
Type=simple
StandardOutput=journal+console
StandardError=journal+console
Restart=always
StartLimitInterval=0
RestartSec=60
EnvironmentFile=/opt/mesosphere/etc/mesos-slave-common
EnvironmentFile=/opt/mesosphere/etc/mesos-slave-public
EnvironmentFile=-/var/lib/dcos/mesos-slave-common
ExecStartPre=/usr/bin/env mkdir -p /var/lib/dcos
ExecStart=/opt/mesosphere/bin/make_disk_resources.py /var/lib/dcos/mesos-resources
