# packages/dcos-diagnostics/extra/dcos-diagnostics-mesos-state.timer
[Unit]
Description=Timer for triggering Mesos state

[Timer]
OnCalendar=daily
RandomizedDelaySec=3600