# packages/dcos-ui-update-service/extra/dcos-ui-update-service.socket
[Unit]
Description=DC/OS UI Update Service Socket: socket for the dcos ui update service
[Socket]
ListenStream=/run/dcos/dcos-ui-update-service.sock
SocketUser=root
SocketGroup=dcos_adminrouter
SocketMode=0660
