# packages/erlang/extra/dcos-epmd.service
[Unit]
Description=Erlang Port Mapping Daemon (EPMD): facilitates communication between distributed Erlang programs

[Service]
User=dcos_epmd
Restart=always
StartLimitInterval=0
RestartSec=5
LimitNOFILE=16384
WorkingDirectory=${PKG_PATH}
EnvironmentFile=/opt/mesosphere/etc/epmd.env
EnvironmentFile=/opt/mesosphere/environment
ExecStart=${PKG_PATH}/bin/epmd -port ${ERL_EPMD_PORT} -relaxed_command_check
Environment=HOME=/opt/mesosphere
