# packages/pkgpanda-api/extra/dcos-pkgpanda-api.socket
[Unit]
Description=DC/OS Component Package Manager (Pkgpanda) Socket: socket for DC/OS Component Package Manager

[Socket]
ListenStream=/run/dcos/pkgpanda-api.sock
SocketUser=root
SocketGroup=dcos_adminrouter
SocketMode=0660
