# packages/pkgpanda-api/extra/dcos-pkgpanda-api.socket
[Unit]
Description=Pkgpanda API Socket: socket for Pkgpanda API

[Socket]
ListenStream=/run/dcos/pkgpanda-api.sock
